# docker/osh-test-cache/etc/supervisor/conf.d/cache.conf
[supervisord]
nodaemon=true
user=root

[program:redis]
command=/usr/bin/redis-server /etc/redis/redis.conf
user=root
autostart=true
autorestart=true
stopsignal=QUIT
redirect_stderr=true
