# docker/test-web/etc/supervisor/conf.d/web.conf
[supervisord]
nodaemon=true
user = root

[program:fpm]
command = /usr/sbin/php-fpm7.0 --fpm-config /etc/php/7.0/fpm/php-fpm.conf
user = root
autostart = true
autorestart = true
redirect_stderr = true

[program:nginx]
command=/usr/sbin/nginx
user = root
autostart = true
autorestart = true
redirect_stderr = true